# debian/decades-listener/DEBIAN/postrm
#!/bin/sh
if [ "$1" = purge ]; then
        update-rc.d twisted-decades-listener remove >/dev/null
fi
