# debian/decades-listener/DEBIAN/postinst
#!/bin/sh
update-rc.d twisted-decades-listener defaults >/dev/null
invoke-rc.d twisted-decades-listener start

