# debian/postinst
#!/bin/sh
update-rc.d decades defaults >/dev/null
#DEBHELPER#
