# debian/postinst
#!/bin/sh
update-rc.d decades-listener defaults >/dev/null
invoke-rc.d decades-listener start
#DEBHELPER#
