# debian/prerm
#!/bin/sh
invoke-rc.d decades-listener stop
#DEBHELPER#
