# debian/prerm
#!/bin/sh
invoke-rc.d twisted-decades-listener stop
#DEBHELPER#
