# debian/decades/DEBIAN/postrm
#!/bin/sh
if [ "$1" = purge ]; then
        update-rc.d decades remove >/dev/null
fi
