# debian/decades/DEBIAN/postinst
#!/bin/sh
update-rc.d decades defaults >/dev/null
invoke-rc.d decades start

